# systemd-user/pwmfreq.timer
[Unit]
Description=Get PWM data

[Timer]
OnBootSec=5sec
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
