# systemd-user/pwmfreq.service
[Unit]
Description=Get PWM data to tmpfs

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '%h/bin/pwmfreq read > /tmp/pwmfreq 2>/dev/null'
