# systemd-user/redshift.service
[Unit]
Description=Set display colortemp

[Service]
ExecStart=/usr/bin/redshift
Restart=on-failure

[Install]
WantedBy=default.target
